# packaging/systemd/helmd.service
[Unit]
Description=helmd — Stream Deck control surface daemon
After=network.target

[Service]
Type=simple
ExecStart={{PYTHON_EXECUTABLE}} -m helmd
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.config/helm/logs/helmd.log
StandardError=append:%h/.config/helm/logs/helmd.err

[Install]
WantedBy=default.target
